# files/fsw.service.in
[Unit]
Description=Factory Software
Wants=dbus.service
After=dbus.service
ConditionVirtualization=!container
DefaultDependencies=yes

[Service]
Type=simple
ExecStart=@BINARY_PATH@
Restart=on-failure
RestartSec=5

[Install]
WantedBy=sysinit.target
